# root/etc/udev/rules.d/hdmi_sound.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/adabru/setup/bin/hdmi_sound_toggle.sh"
